# lime-api-orange-rpc/files/etc/uci-defaults/99-save-last-known-path-to-internet
#!/bin/sh
echo "saving last known path to the Internet"
echo "*/6 * * * * /usr/sbin/internet_path > /tmp/last_internet_path" >> /etc/crontabs/root

