# lighttpd-ws/files/etc/uci-defaults/disable_uhttpd.sh
#!/bin/sh
echo "Disabling uhttpd, let's use lighttpd instead"
[ -f /etc/init.d/uhttpd ] && /etc/init.d/uhttpd disable
